# templates/jmxtrans.service.epp
[Unit]
Description=JMX Transformer - more than meets the eye
After=syslog.target network.target

[Service]
Type=forking
User=<%= $user %>
# Run ExecStartPre with root-permissions
PermissionsStartOnly=true
RuntimeDirectory=jmxtrans
PIDFile=/var/run/jmxtrans/jmxtrans.pid
ExecStart=/usr/share/jmxtrans/bin/jmxtrans.sh start

[Install]
WantedBy=multi-user.target
